# launchers/systemd/webterm-agent.service
# webterm agent unit (one always-on shell registered with the broker;
# additional terminals come from the broker's POST /launch). Install:
#   sudo cp launchers/systemd/webterm-agent.service /etc/systemd/system/
#   edit User/WorkingDirectory/ExecStart paths, broker URL and token, then
#   sudo systemctl daemon-reload && sudo systemctl enable --now webterm-agent
#
# The token below is REQUIRED (#142): every connection needs one, including
# loopback. It must MATCH the broker's - read it with
#   python3 -m webterm.broker --print-token
# and note that an agent started with the wrong token is closed with 4401 and
# retries every 10s without ever registering.
[Unit]
Description=webterm agent (headless PTY producer)
After=network-online.target webterm-broker.service
Wants=network-online.target

[Service]
Type=simple
User=webterm
WorkingDirectory=/opt/web_terminal
Environment=BROWSERLAND_BROKER_URL=ws://127.0.0.1:4445/browserland
Environment=WEB_TERMINAL_TOKEN=changeme
ExecStart=/opt/web_terminal/launchers/run-agent.sh
# on-failure (not always): a deliberate UI Terminate of the always-on shell
# exits 0 and stays gone; a crash exits non-zero and respawns. Boot start is
# unaffected. Note: typing `exit` in the shell (code 0) also no longer respawns.
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
